# files/kni.service
[Service]
Restart=always
ExecStart=
ExecStart=/usr/bin/network-runtime
[Install]
WantedBy=multi-user.target
